# debian/xymon.preinst
#! /bin/sh
# preinst script for xymon

set -e

case "$1" in
    install|upgrade)
        # stop the client when the server is installed
        invoke-rc.d xymon-client stop || true
        ;;
esac

#DEBHELPER#

exit 0
